# crontab
* * * * * python /main.py >> /dev/console